# scripts/crontab
2 * * * * /bin/sh /home/aluxious/mamtask.sh
